# trading-bot.service
[Unit]
Description=AI Trading Bot - Next.js Production Server
After=network.target

[Service]
Type=simple
User=aiminer2
WorkingDirectory=/home/aiminer2/Trading Project/deepseek-ui

# Build first, then start the production server.
# ExecStartPre runs blocking — service only starts after build succeeds.
ExecStartPre=/usr/bin/npm run build
ExecStart=/usr/bin/npm start

Restart=always
RestartSec=10
StandardOutput=append:/home/aiminer2/Trading Project/trading-bot.log
StandardError=append:/home/aiminer2/Trading Project/trading-bot-error.log

# Environment variables
Environment="NODE_ENV=production"
Environment="PATH=/usr/bin:/usr/local/bin:/home/aiminer2/.nvm/versions/node/v20.11.0/bin"

[Install]
WantedBy=multi-user.target
